# deploy/systemd/genie-core.service
[Unit]
Description=GeniePod Core — Local Home AI Engine
Documentation=https://github.com/GeniePod/genie-claw
# Order after both LLM units so genie-core never starts before whichever
# backend is configured. genie-ai-runtime.service is the Jetson default;
# genie-llm.service is the llama.cpp fallback. They Conflict= each other,
# so only one will ever be active at a time and the unused After= entry
# is a no-op.
After=network.target genie-mqtt.service genie-ai-runtime.service genie-llm.service
Wants=genie-mqtt.service genie-ai-runtime.service

[Service]
Type=simple
ExecStart=/opt/geniepod/bin/genie-core
WorkingDirectory=/opt/geniepod
Environment=GENIEPOD_CONFIG=/etc/geniepod/geniepod.toml
Environment=RUST_LOG=info
# Stdin must NOT be a terminal, so genie-core runs in HTTP server mode.
StandardInput=null
Restart=always
RestartSec=3

# CPU pinning for voice latency stability (issue #25). Core 5 hosts the
# orchestrator and every audio child it spawns (piper, sox, deep-filter,
# arecord, aplay) — they inherit this affinity. Keeping the TTS pipeline on
# one core shields whisper (2-3) and llama-server (4) from preemption during
# decode. arecord runs on core 5 by inheritance (not 0-1 as the issue's
# system-bucket bullet suggests); revisit with per-spawn taskset only if
# capture jitter shows up under load.
CPUAffinity=5

# Needs GPU (Whisper TRT, Piper TTS), audio (ALSA/I2S), and network.
ProtectSystem=no
SupplementaryGroups=video render audio

[Install]
WantedBy=geniepod.target
